# launchd/com.dexterpratt.agent.rcorona-morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/dexterpratt</string>
		<key>PATH</key>
		<string>/Users/dexterpratt/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>com.dexterpratt.agent.rcorona-morning</string>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/dexterpratt/.local/bin/run-scheduled-agent.sh</string>
		<string>rcorona</string>
		<string>morning</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/Users/dexterpratt/.claude/logs/launchd/rcorona-morning.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/dexterpratt/.claude/logs/launchd/rcorona-morning.out.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>7</integer>
		<key>Minute</key>
		<integer>39</integer>
	</dict>
</dict>
</plist>
