# launchd/com.dexterpratt.agent.rdaneel-quals-prep.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dexterpratt.agent.rdaneel-quals-prep</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/dexterpratt/.local/bin/run-scheduled-agent.sh</string>
        <string>rdaneel</string>
        <string>quals-prep</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Day</key><integer>15</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Day</key><integer>15</integer>
            <key>Hour</key><integer>15</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Day</key><integer>15</integer>
            <key>Hour</key><integer>21</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Day</key><integer>16</integer>
            <key>Hour</key><integer>3</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Day</key><integer>16</integer>
            <key>Hour</key><integer>9</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Day</key><integer>16</integer>
            <key>Hour</key><integer>15</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/dexterpratt/.claude/logs/launchd/rdaneel-quals-prep.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dexterpratt/.claude/logs/launchd/rdaneel-quals-prep.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/dexterpratt/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/dexterpratt</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
